# local-ai-assistant/scripts/local-ai-assistant.service
[Unit]
Description=Local AI Assistant
After=network.target ollama.service

[Service]
Type=simple
User=axmh
WorkingDirectory=/home/axmh/Desktop/AxisAIAssistant/local-ai-assistant
ExecStart=/home/axmh/Desktop/AxisAIAssistant/local-ai-assistant/.venv/bin/python /home/axmh/Desktop/AxisAIAssistant/local-ai-assistant/main.py
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/axmh/.Xauthority

[Install]
WantedBy=default.target
